# infra/tmp-fakehwclock.mount
[Unit]

[Mount]
What=tmpfs
Where=/tmp/fakehwclock
Type=tmpfs

[Install]
WantedBy=multi-user.target
