# infra/var-tmp-dietpi-logs.mount
[Unit]

[Mount]
What=tmpfs
Where=/var/tmp/dietpi/logs
Type=tmpfs

[Install]
WantedBy=multi-user.target
