# infra/var-lib-dhcp.mount
[Unit]

[Mount]
What=tmpfs
Where=/var/lib/dhcp
Type=tmpfs

[Install]
WantedBy=multi-user.target
